# systemd/dirigera-firmware-check.timer
[Unit]
Description=IKEA DIRIGERA Firmware Check alle 6h

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
